# root/defaults/update-nessus.cron
1 1 * * * /command/with-contenv bash /defaults/update-nessus.sh >/dev/null 2>&1
