# root/defaults/update-nessus.cron
1 1 * * * /bin/bash /defaults/update-nessus.sh >/dev/null 2>&1
