# root/etc/services.d/cron/run
#!/usr/bin/with-contenv bash

crond -n -s