# root/etc/services.d/nessus/run
#!/command/with-contenv bash

# NOTE - TEMPORARILY DISABLE NOROOT WHILE TROUBLESHOOTING SETCAP ISSUES!
#s6-setuidgid abc /opt/nessus/sbin/nessus-service --no-root
/opt/nessus/sbin/nessus-service